# centos7.mod_wsgi-express/roles/supervisor/templates/supervisor/conf.d/example/bdai_exampleapi-4.ini
[program:bdai_exampleapi-4]
command={{mod_wsgi_express}} start-server 
    --server-root {{wsgi_app_location}}/ws-4/bdai_exampleapi 
    --host 0.0.0.0
    --port 9004
    --user apache
    --processes 1
    --threads 100
    app.py
directory={{wsgi_app_location}}/ws-4/bdai_exampleapi
environment=prometheus_multiproc_dir="{{prometheus_data_location}}/ws-4/bdai_exampleapi"
