# centos7.mod_wsgi-express/roles/mod_wsgi-express/files/supervisor/conf.d/bdai_exampleapi-1.ini
[program:bdai_exampleapi-1]
command=/usr/local/bin/mod_wsgi-express start-server
    --server-root /data/www/wsgi/ws-1/bdai_exampleapi
    --host 0.0.0.0
    --port 9001
    --user apache
    --processes 1
    --threads 100
    app.py
directory=/data/www/wsgi/ws-1/bdai_exampleapi
environment=prometheus_multiproc_dir="/data/www/prometheus/ws-1"
