# centos7.mod_wsgi-express/roles/mod_wsgi-express/files/supervisor/conf.d/bdai_recommendapi-1.ini
[program:bdai_recommendapi-1]
command=/usr/local/bin/mod_wsgi-express start-server --server-root /data/www/wsgi.run/ws-1/bdai_recommendapi --host 0.0.0.0 --port 9001 --user apache app.py
directory=/data/www/wsgi/ws-1/bdai_recommendapi
environment=prometheus_multiproc_dir="/data/www/prometheus/ws-1"
