# centos7.mod_wsgi-express/roles/supervisor/templates/supervisord.service
[Unit]
Description=Supervisor Daemon
After=network.target network-online.target

[Service]
Type=simple
ExecStart={{supervisord}} -c /etc/supervisor/supervisord.conf
KillMode=process
KillSignal=SIGINT
Restart=on-failure

[Install]
WantedBy=multi-user.target
