# centos7.mod_wsgi-express/roles/mod_wsgi-express/files/supervisor/conf.d/example/bdai_exampleapi-4.ini
[program:bdai_recommendapi-4]
command=/usr/local/bin/mod_wsgi-express start-server --server-root /data/www/wsgi.run/ws-4/bdai_recommendapi --host 0.0.0.0 --port 9004 --user apache app.py
directory=/data/www/wsgi/ws-4/bdai_recommendapi
environment=prometheus_multiproc_dir="/data/www/prometheus/ws-4"
